# .config/autostart/walchange.desktop
[Desktop Entry]
Type=Application
Name=walchange
Exec=$HOME/.local/bin/walchange

